# .config/systemd/user/dunst.service
[Unit]
Description=dunst

[Service]
Type=simple
ExecStart=/usr/bin/env dunst 
ExecStop=/usr/bin/pkill dunst

[Install]
WantedBy=wm.target
