# .config/systemd/user/jetbrains-toolbox.service
[Unit]
Description=jetbrains-toolbox
Wants=taffybar.service
After=taffybar.service

[Service]
ExecStart=/usr/bin/jetbrains-toolbox --minimize
Restart=always

[Install]
WantedBy=wm.target
