# .config/systemd/user/status-notifier-watcher.service
[Unit]
Description=status-notifier-watcher

[Service]
Type=simple
ExecStart=/usr/bin/env status-notifier-watcher
ExecStop=/usr/bin/pkill -f status-notifier-watcher

[Install]
WantedBy=wm.target
