# .config/systemd/user/status-notifier-watcher.service
[Unit]
Description=status-notifier-watcher

[Service]
Type=simple
ExecStart=/usr/bin/env /home/patrick/.asdf/shims/status-notifier-watcher
ExecStop=/usr/bin/pkill status-notifier-watcher

[Install]
WantedBy=wm.target
