# .config/systemd/user/wallpaper.service
[Unit]
Description=wallpaper

[Service]
Type=oneshot
ExecStart=/usr/bin/env nitrogen --restore

[Install]
WantedBy=wm.target
